# uploads/backup.timer
[Unit]
Description=Run backup script at 3 AM WIB
After=network.target

[Timer]
OnCalendar=03:00:00
Persistent=true
Unit=backup.service

[Install]
WantedBy=timers.target
